# home/conf/supervisord.conf
[program:home-dashboard]
command = /home/pi/.virtualenvs/home/bin/home dashboard
user = pi
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/home-dashboard.log
stderr_logfile = /var/log/supervisor/home-dashboard-error.log
environment = HOME_SETTINGS=/home/pi/home.py

[program:home-collect]
command = /home/pi/.virtualenvs/home/bin/home collect
user = pi
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/home-collect.log
stderr_logfile = /var/log/supervisor/home-collect-error.log
environment = HOME_SETTINGS=/home/pi/home.py
