# home/conf/supervisord.conf
[program:home-dashboard]
command = /home/pi/.virtualenvs/home/bin/uwsgi -s /tmp/uwsgi.home.sock --chmod-socket=666 -w home.__main__:app
directory = /home/pi/.virtualenvs/home/lib/python3.4/site-packages
user = pi
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/home-dashboard.log
stderr_logfile = /var/log/supervisor/home-dashboard-error.log
environment=HOME_SETTINGS="/home/pi/home_settings.py"

[program:home-collect]
command = /home/pi/.virtualenvs/home/bin/home collect
directory = /home/pi/.virtualenvs/home/lib/python3.4/site-packages
user = pi
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/home-collect.log
stderr_logfile = /var/log/supervisor/home-collect-error.log
environment=HOME_SETTINGS="/home/pi/home_settings.py"
